# deploy/systemd/surveys.service
[Unit]
Description=Dostoevsky Soul Trial
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/surveys
Environment=NODE_ENV=production
Environment=OPEN_ACCESS=1
Environment=HOST=127.0.0.1
Environment=PORT=8787
Environment=SESSION_SECRET=replace-with-a-long-random-secret
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
